# vps-setup/agent-template/systemd/brief-summary.service.tmpl
[Unit]
Description={{TENANT_PERSON_FIRST_NAME}} Overview brief-summary (1-3 sentence chief-of-staff strip)
Documentation=https://github.com/newmindsgroup/{{TENANT_BRAND_REPO_NAME}}
After=network.target

[Service]
Type=oneshot
User={{TENANT_LINUX_USER}}
Group={{TENANT_LINUX_GROUP}}
WorkingDirectory={{TENANT_AGENT_HOME}}
Environment="HOME={{TENANT_USER_HOME}}"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:{{TENANT_USER_HOME}}/.bun/bin:{{TENANT_USER_HOME}}/.local/bin"
Environment="TENANT_AGENT_HOME={{TENANT_AGENT_HOME}}"
ExecStart=/usr/bin/python3 {{TENANT_AGENT_HOME}}/scripts/brief-summary.py
StandardOutput=append:{{TENANT_AGENT_HOME}}/logs/brief-summary.log
StandardError=append:{{TENANT_AGENT_HOME}}/logs/brief-summary.log
TimeoutStartSec=90
