# vps-setup/agent-template/systemd/disk-space-watcher.service.tmpl
[Unit]
Description={{TENANT_PERSON_FIRST_NAME}} Gonell disk-space watcher (proactive Telegram nudges)
Documentation=https://github.com/newmindsgroup/{{TENANT_BRAND_REPO_NAME}}
After=network.target

[Service]
Type=oneshot
User={{TENANT_LINUX_USER}}
Group={{TENANT_LINUX_GROUP}}
WorkingDirectory={{TENANT_AGENT_HOME}}
Environment="HOME={{TENANT_USER_HOME}}"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:{{TENANT_USER_HOME}}/.bun/bin:{{TENANT_USER_HOME}}/.local/bin"
Environment="TENANT_AGENT_HOME={{TENANT_AGENT_HOME}}"
ExecStart={{TENANT_AGENT_HOME}}/scripts/disk-space-watcher.sh
StandardOutput=append:{{TENANT_AGENT_HOME}}/logs/disk-space-watcher.log
StandardError=append:{{TENANT_AGENT_HOME}}/logs/disk-space-watcher.log
TimeoutStartSec=60
