# vps-setup/agent-template/systemd/rules-engine.service.tmpl
[Unit]
Description=Behavioral rules engine — evaluates rules/*.yaml against /api/* state
After=network.target dashboard-chat.service
Wants=network.target

[Service]
Type=oneshot
User={{TENANT_LINUX_USER}}
Group={{TENANT_LINUX_GROUP}}
WorkingDirectory={{TENANT_AGENT_HOME}}
Environment="HOME={{TENANT_USER_HOME}}"
Environment="PATH=/usr/local/bin:/usr/bin:/bin:{{TENANT_USER_HOME}}/.local/bin"
ExecStart=/usr/bin/python3 {{TENANT_AGENT_HOME}}/scripts/rules-engine.py
StandardOutput=append:{{TENANT_AGENT_HOME}}/logs/rules-engine.log
StandardError=append:{{TENANT_AGENT_HOME}}/logs/rules-engine.log
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
